# certrenew.target
[Unit]
Description=Services which need reloaded with SSL certs are updated
PropagatesReloadTo=nginx postfix dovecot
